# cherrypy/systemd/my_sensors.service
[Unit]
Description=Grab Sensor Data

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/proxmox-sensors/my_sensors.py

