# cherrypy/systemd/my_sensors.timer
[Unit]
Description=Grab Sensor Data Timer

[Timer]
OnUnitActiveSec=3000s
OnBootSec=30s

[Install]
WantedBy=timers.target
